# service/smartdoor_backup.timer
[Unit]
Description = Backup timer for Smartdoor Host database

[Timer]
OnCalendar = weekly
Persistent = true

[Install]
WantedBy = timers.target
