# service/smartdoor_backup.timer
[Unit]
Description = Backup timer for Smartdoor Host database

[Service]
OnCalendar = weekly
Persistent = true

[Install]
WantedBy = timers.target
